# code/root/etc/s6-overlay/s6-rc.d/code-installer/run
#!/usr/bin/env sh

if [ "$(id -u)" -eq 1000 ]; then
  echo USER OK
  echo Installing package in VENV
  . /config/venv/bin activate
  while true; do
    echo trying to install package
    pyproject=$(find /config/workspace -name pyproject.toml)
    if [ -f "${pyproject}" ]; then
      echo trying to install ${pyproject}
      cd $(dirname ${pyproject})
      /config/venv/bin/pip install -e .[lint,test] && break
    else
      still no project found
    fi
    sleep 10
    echo failed, retrying
  done
else 
  echo switching USER
  sudo -u abc $0 $*
fi

