# code/root/etc/s6-overlay/s6-rc.d/code-installer/run
#!/usr/bin/env sh

echo run CALLED

while true; do
  pip install -e .[lint|test] && break
  sleep 10
done
chmod -R a+rwX *egg-info
[ -e .venv ] || ln -s ../venv .venv


